# debian/nntp-rs.nntp-rs.service
[Unit]
Description=NNTP Server
After=network.target
ConditionPathExists=/etc/news/config.cfg

[Service]
LimitMEMLOCK=infinity
LimitNOFILE=65536
ExecStart=/usr/sbin/nntp-rs-server -c /etc/news/config.cfg serve
KillMode=process

[Install]
WantedBy=multi-user.target

